# configs/systemd/mgs-wal.service
[Unit]
Description=MGS adaptive colors — inotify watcher
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mgs-wal-watch
Restart=on-failure
RestartSec=3
Environment=PATH=/usr/bin:/usr/local/bin:%h/.local/bin
PassEnvironment=HYPRLAND_INSTANCE_SIGNATURE WAYLAND_DISPLAY XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS

[Install]
WantedBy=graphical-session.target
